# WordPress/nmu-rs-cloudfiles/cloudfiles-worker.conf
[program:cloudfiles-worker]
command=/var/www/html/build/wp-content/plugins/nmu-rs-cloudfiles/cloudfiles-worker.php
autostart=true